# 《DevOps和自动化运维实践》第6章/srv/salt/base/saltcheck/salt_agent.cron
*/5 * * * * root /usr/bin/saltcheck.sh >>/var/log/salt.log
#即每隔5分钟就执行/usr/bin/saltecheck.sh命令来检查下salt-minion服务
